# mppsolar-docker/root/etc/s6-overlay/s6-rc.d/mpp-install/up
#!/bin/execlineb -P
/bin/bash /etc/s6-overlay/scripts/mpp-install.sh
