# ha-mppsolar/rootfs/etc/services.d/web-ui/run
#!/bin/sh
for f in /run/s6/container_environment/*; do
    [ -f "$f" ] && export "$(basename "$f")=$(cat "$f")"
done
exec python3 /opt/web-ui/app.py
